# dev-stack/claude-logger.service
[Unit]
Description=Claude Code Background Logger
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=adam
WorkingDirectory=/home/adam/Desktop/project-bolt-sb1-1c6erknt/ATS-resume-wesite/dev-stack
ExecStart=/usr/bin/python3 /home/adam/Desktop/project-bolt-sb1-1c6erknt/ATS-resume-wesite/dev-stack/claude_background_logger.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-logger

[Install]
WantedBy=multi-user.target